# deploy/kitan-telegram-self-improve.service
[Unit]
Description=Kitan Telegram Self-Improve Nudge
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=foggen
WorkingDirectory=/home/foggen/kitan-telegram
EnvironmentFile=-/home/foggen/kitan-telegram/.env
ExecStart=/home/foggen/kitan-telegram/scripts/self_improve_nudge.sh

