# deploy/kitan-telegram-self-test.timer
[Unit]
Description=Run Kitan Telegram reliability self-test periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=kitan-telegram-self-test.service

[Install]
WantedBy=timers.target

